# data/scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:watchtopdf]
command=/var/www/html/scripts/watchtopdf.sh
autostart=true
autorestart=true
stopwaitsecs=10

[program:topdf]
command=/var/www/html/scripts/topdf.sh
autostart=true
autorestart=true
stopwaitsecs=10

